# supervisor.conf
[program:nginx]
directory=/app/leanEngine_app
command=gunicorn -w 5 -k gevent wsgi:application
